# etc/service/redis-server/run
#!/bin/sh

set -e
chown -R redis:redis /var/lib/redis
exec redis-server
